# darkstar-ahbroker/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:broker]
directory=/darkstar/bin
command=/usr/bin/python3 ./broker.py
stdout_events_enabled=true
stderr_events_enabled=true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 5
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler